# templates/systemd/kubexm-cert-renew-kubeadm.service.tmpl
[Unit]
Description=Kubernetes Certificate Renewal Check (kubeadm)
Documentation=https://kubernetes.io/docs/reference/setup-tools/kubeadm/kubeadm-certs/
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# kubeadm 模式使用 kubeadm certs renew all
ExecStart=/bin/bash -c '\
  DAYS_BEFORE=${CERT_RENEW_DAYS_BEFORE:-30}; \
  for cert in /etc/kubernetes/pki/*.crt; do \
    if [ -f "$cert" ]; then \
      end_date=$(openssl x509 -enddate -noout -in "$cert" 2>/dev/null | cut -d= -f2); \
      end_ts=$(date -d "$end_date" +%s 2>/dev/null); \
      now_ts=$(date +%s); \
      days_left=$(( (end_ts - now_ts) / 86400 )); \
      if [ "$days_left" -lt "$DAYS_BEFORE" ]; then \
        echo "Certificate $cert expires in $days_left days, renewing..."; \
        /usr/bin/kubeadm certs renew all; \
        systemctl restart kubelet; \
        exit 0; \
      fi; \
    fi; \
  done; \
  echo "All certificates are valid for more than $DAYS_BEFORE days"'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
